# deployment/systemd/django-q.service
[Unit]
Description=Django-Q Cluster for Embassy PWA
After=network.target postgresql.service

[Service]
Type=simple
User=webapp
Group=www-data
WorkingDirectory=/home/webapp/embassy/backend
Environment="PATH=/home/webapp/embassy/venv/bin"
EnvironmentFile=/home/webapp/embassy/backend/.env
ExecStart=/home/webapp/embassy/venv/bin/python manage.py qcluster
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

